# home/.config/systemd/user/polybar-update.path
[Unit]
Description=Watch apt updates file

[Path]
PathModified=/var/cache/machinepwn/updates.txt

[Install]
WantedBy=default.target
