# home/.config/systemd/user/polybar-update.service
[Unit]
Description=Update polybar updates hook

[Service]
ExecStart=/usr/bin/polybar-msg action updates hook 0
